# lb-agent/etc/supervisord.d/lb-agent.ini
[program:lb-agent]
command=/usr/bin/python /usr/local/lb-agent/lb-agent.py
autostart=true
autorestart=true
user=root
